# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable wiser-exporter.service
systemctl restart wiser-exporter.service
