# packaging/scripts/preremove.sh
systemctl stop wiser-exporter.service
systemctl disable wiser-exporter.service
